# ci-cd/crontab
# Delete old docker-containers, -volumes, -networks and -images
0 5 * * 1 /usr/bin/docker system prune --all --force
# letzte Zeile muss kommentarzeile oder leer sein
